# .config/systemd/user/wp-slideshow.service
[Unit]
Description=Wallpaper Slideshow for Hyprland

[Service]
Type=oneshot
ExecStart=%h/.config/bin/wp-slideshow.sh
