# .config/systemd/user/wp-slideshow.service
[Unit]
Description=Wallpaper Slideshow for Hyprland
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/bin/wp-slideshow.sh
